# devpi-server/supervisord.conf
[supervisord]
nodaemon=true

[supervisorctl]

[program:devpi-server]
command=/usr/local/bin/devpi-server -c /devpi-server.yaml
priority = 999
startsecs = 5
redirect_stderr = True
autostart = True

